# privetd/init/privetd.conf
# Copyright 2014 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Privet protocol handler daemon"
author        "chromium-os-dev@chromium.org"

start on starting system-services and stopped iptables and stopped ip6tables
stop on stopping system-services

env PRIVETD_LOG_LEVEL=0
env PRIVETD_HTTP_PORT=8080
env PRIVETD_HTTPS_PORT=8081
env PRIVETD_ENABLE_PING=false
env PRIVETD_DISABLE_SECURITY=false
env PRIVETD_STATE_PATH=
env PRIVETD_CONFIG_PATH=
env PRIVETD_DEVICE_WHITELIST=
# TODO(zqiu): Temporary
env APMANAGER_DHCP_SERVER_PORT=67

pre-start script
  # Add a rule to the firewall to allow HTTP traffic except from the
  # loopback interface (to prevent e.g. Chrome from connecting.)
  # TODO(avakulenko): localhost is left enabled for now for testing...
  #iptables -I INPUT -i lo -p tcp --dport ${PRIVETD_HTTP_PORT} -j REJECT
  #ip6tables -I INPUT -i lo -p tcp --dport ${PRIVETD_HTTP_PORT} -j REJECT
  #iptables -I INPUT -i lo -p tcp --dport ${PRIVETD_HTTPS_PORT} -j REJECT
  #ip6tables -I INPUT -i lo -p tcp --dport ${PRIVETD_HTTPS_PORT} -j REJECT
  iptables -A INPUT -p tcp --dport ${PRIVETD_HTTP_PORT} -j ACCEPT
  ip6tables -A INPUT -p tcp --dport ${PRIVETD_HTTP_PORT} -j ACCEPT
  iptables -A INPUT -p tcp --dport ${PRIVETD_HTTPS_PORT} -j ACCEPT
  ip6tables -A INPUT -p tcp --dport ${PRIVETD_HTTPS_PORT} -j ACCEPT

  # TODO(zqiu): Temporary hack for adding firewall rule to accept DHCP traffic
  # on behalf of apmanager. Putting it here instead of in the apmanager upstart
  # script to avoid a race condition when multiple upstart scripts tried to
  # execute iptables command at the same time. Will delete this when the
  # support for runtime firewall hole punching in permission_broker is
  # completed (crbug.com/435400).
  iptables -I INPUT -p udp --dport ${APMANAGER_DHCP_SERVER_PORT} -j ACCEPT
  iptables -I INPUT -i lo -p udp --dport ${APMANAGER_DHCP_SERVER_PORT} -j REJECT

  # privetd runs as its own user, and can't create files on the filesystem.
  # Create appropriate workspaces accordingly.
  mkdir -m 0755 -p /var/lib/privetd
  chown -R privetd:privetd /var/lib/privetd
end script

# Minijail actually forks off our desired process
expect fork

exec minijail0 -i -g privetd -u privetd /usr/bin/privetd \
    --v="${PRIVETD_LOG_LEVEL}" \
    --http_port="${PRIVETD_HTTP_PORT}" \
    --https_port="${PRIVETD_HTTPS_PORT}" \
    --disable_security="${PRIVETD_DISABLE_SECURITY}" \
    --enable_ping="${PRIVETD_ENABLE_PING}" \
    --config_path="${PRIVETD_CONFIG_PATH}" \
    --state_path="${PRIVETD_STATE_PATH}" \
    --device_whitelist="${PRIVETD_DEVICE_WHITELIST}"

post-stop script
  # Delete the rules we previously added
  iptables -D INPUT -p tcp --dport ${PRIVETD_HTTP_PORT} -j ACCEPT
  ip6tables -D INPUT -p tcp --dport ${PRIVETD_HTTP_PORT} -j ACCEPT
  iptables -D INPUT -p tcp --dport ${PRIVETD_HTTPS_PORT} -j ACCEPT
  ip6tables -D INPUT -p tcp --dport ${PRIVETD_HTTPS_PORT} -j ACCEPT

  # TODO(zqiu): refer to comments above for adding these rules.
  iptables -D INPUT -p udp --dport ${APMANAGER_DHCP_SERVER_PORT} -j ACCEPT
  iptables -D INPUT -i lo -p udp --dport ${APMANAGER_DHCP_SERVER_PORT} -j REJECT

  # TODO(avakulenko): uncomment below
  #iptables -D INPUT -i lo -p tcp --dport ${PRIVETD_HTTP_PORT} -j REJECT
  #ip6tables -D INPUT -i lo -p tcp --dport ${PRIVETD_HTTP_PORT} -j REJECT
  #iptables -D INPUT -i lo -p tcp --dport ${PRIVETD_HTTPS_PORT} -j REJECT
  #ip6tables -D INPUT -i lo -p tcp --dport ${PRIVETD_HTTPS_PORT} -j REJECT
end script
