# brdebug/udev/99-setup-debug-link.rules
ACTION=="add|change", SUBSYSTEM=="usb", ATTRS{idVendor}=="0b95", ATTRS{idProduct}=="7e2b", RUN+="/usr/sbin/setup-debug-link"
