# brdebug/udev/99-setup-veth-link.rules
ACTION=="add|change", SUBSYSTEM=="usb", ATTRS{idVendor}=="13b1", ATTRS{idProduct}=="0041", RUN+="/usr/sbin/setup-veth-link"
