# wimax_manager/wimax_manager.conf
# Copyright (c) 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# wimax_manager upstart job

# Minimum log level defined in base/logging.h.
#   0:INFO, 1:WARNING, 2:ERROR, 3:ERROR_REPORT, 4:FATAL
# Set to log only WARNING or above by default.
env WIMAX_MANAGER_MINLOGLEVEL=1

start on starting system-services
stop on stopping system-services

respawn
expect fork

exec wimax-manager -minloglevel="${WIMAX_MANAGER_MINLOGLEVEL}"
