# init/test-init/autoreboot.conf
# Copyright (c) 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Enables power on after power failure for DUT's in the lab"
author          "chromium-os-dev@chromium.org"

# For test images being ran on Intel-based devices in the lab, set the
# correct pci bit so that it will reboot on a power cycle.

start on starting boot-services
script
  LAB_FILE="/mnt/stateful_partition/.labmachine"
  if [ -f "$LAB_FILE" ]; then
    # Future devices which we can control power on after power failure
    # can be added below.
    if grep -q CHROMEOS_RELEASE_BOARD=stumpy /etc/lsb-release; then
      /usr/sbin/setpci -s 00:1f.0 0xA4.b=4
    fi
  fi
end script