# init/test-init/factory-init.conf
description    "Initialize factory test environment"
author         "chromium-os-dev@chromium.org"

# For factory test to perform early initialization Ex: changing UI start
# parameters, disable or add upstart jobs...
# Hint: ui.conf starts on "started boot-services", so this factory-init should
# always run before UI starts.

start on starting boot-services
task

script
  . /usr/share/cros/factory_utils.sh
  if is_factory_test_mode; then
    factory_init="/usr/local/factory/init/startup"
    [ -x "${factory_init}" ] && "${factory_init}"
  fi
end script
