# init/test-init/factorylog.conf
# Copyright (c) 2013 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Print Chrome OS factory log to tty3"
author        "chromium-os-dev@chromium.org"

start on started factory
stop on stopping factory

respawn

pre-start script
  . /usr/share/cros/factory_utils.sh
  if ! is_factory_mode; then
    stop
  fi
end script

script
  tail -n 48 -F /var/factory/log/factory.log >/dev/tty3
end script
