# init/test-init/factory-splash.conf
# Copyright (c) 2013 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description    "Overlay factory test logo onto boot splash"
author         "chromium-os-dev@chromium.org"

start on stopping startup and stopping boot-splash
task

script
  . /usr/share/cros/factory_utils.sh
  if is_factory_mode; then
    ply-image "/usr/local/factory/misc/boot_splash.png" --offset=0,+80
  fi
end script
