# init/test-init/inhibit-cros-disks.conf
# Copyright (c) 2014 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Inhibit cros-disks for factory testing"
author        "chromium-os-dev@chromium.org"

start on starting cros-disks
stop on stopped cros-disks

# To start cros-disks regardless of the presence of factory enabled
# tag, do:
#   $ start inhibit-cros-disks disable_inhibit=1
#   $ start cros-disks
env disable_inhibit=0

pre-start script
  . /usr/share/cros/factory_utils.sh
  inhibit_if_factory_mode cros-disks
end script
