# init/test-init/inhibit-laptop-mode-boot.conf
# Copyright (c) 2013 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Inhibit laptop-mode-boot for factory testing"
author        "chromium-os-dev@chromium.org"

start on starting laptop-mode-boot
stop on stopped laptop-mode-boot

# To start laptop-mode-boot regardless of the presence of factory enabled
# tag, do:
#   $ start inhibit-laptop-mode-boot disable_inhibit=1
#   $ start laptop-mode-boot
env disable_inhibit=0

pre-start script
  . /usr/share/cros/factory_utils.sh
  inhibit_if_factory_mode laptop-mode-boot
end script
