# init/test-init/inhibit-laptop-mode-resume.conf
# Copyright (c) 2013 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Inhibit laptop-mode-resume for factory testing"
author        "chromium-os-dev@chromium.org"

start on starting laptop-mode-resume
stop on stopped laptop-mode-resume

# To start laptop-mode-resume regardless of the presence of factory enabled
# tag, do:
#   $ start inhibit-laptop-mode-resume disable_inhibit=1
#   $ start laptop-mode-resume
env disable_inhibit=0

pre-start script
  . /usr/share/cros/factory_utils.sh
  inhibit_if_factory_mode laptop-mode-resume
end script
