# init/test-init/upstart-logging.conf
# Copyright 2016 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Enables verbose Upstart logging for test images."
author          "chromium-os-dev@chromium.org"

# Ideally this should be run as early as possible.
start on startup
task

exec initctl log-priority info
