# power_manager/99-power-supply.rules
ACTION=="add|remove|change", SUBSYSTEM=="power_supply", RUN+="/usr/bin/power_supply_changed"
