# init/acpid.conf
# Copyright (c) 2011 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "ACPI daemon service"
author        "chromium-os-dev@chromium.org"

# acpid is responsible for handling changes in power state, such as
# changes in whether we're running on AC power or battery.

start on starting system-services
stop on stopping system-services

respawn
expect fork
exec acpid --confdir /etc/acpi/events --socketfile /var/run/acpid.socket
