# init/alsa-default.conf
# Copyright (c) 2011 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Set ALSA default mixer settings"
author          "chromium-os-dev@chromium.org"

# On every user session, ensure the microphone is on and set at a
# reasonable level.
start on start-user-session

script
  f=/usr/share/alsa/init/chromeos-per-session
  if [ -f ${f} ]; then
    alsactl init -i ${f}
  else
    amixer sset 'Capture',0 '80%' cap
  fi
end script
