# init/alsa-default.conf
# Copyright (c) 2011 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "Set ALSA default mixer settings"
author "chromium-os-dev@chromium.org"

# When starting a session, ensure the microphone is on and set at a
# reasonable level.

start on start-user-session
task

script
  /usr/bin/amixer sset 'Capture',0 '80%' cap
end script

