# init/audiopolicy.conf
# Copyright (c) 2011 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# audiopolicy

description     "Trigger audio channel switches"
author          "chromium-os-dev@chromium.org"

# This job controls when the headphone jack should enable or disbable.
# It starts roughtly at login time, and is also executed at factory time.
# This script covers both initial config, and a daemonized script
# to poll for headphone hotplugs.

start on starting system-services
stop on stopping system-services

# Located in /usr/sbin/headphone-jack-monitor
exec headphone-jack-monitor
