# init/boot-complete.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# boot-complete
#
# Processing required during boot after login prompt is displayed.
# Primarily, this means calculating boot time statistics.


# Note, this line is modified by chromeos-factoryinstall.ebuild
# and mod_for_factory. Please do not change without also changing
# those references.
start on login-prompt-ready

task

script
  # The login-prompt-ready signal is emitted every time the login
  # screen is displayed so the check below makes sure boot stats are
  # generated only once per boot.
  if [ -n "$(bootstat_last login-prompt-ready)" ]; then
    exit 0
  fi

  bootstat login-prompt-ready

  # If ureadahead is still running, then presumably it's performing
  # a boot trace.  We only care about tracing until login prompt, so
  # stop the process now.
  #
  # We have to ignore the return status in case ureadahead has
  # already completed, because pkill helpfully exits with a non-zero
  # status when no process was found.
  pkill '^ureadahead$' || true

  # Report boot time and disk usage metrics to UMA.
  send_boot_metrics &
end script
