# init/cgroups.conf
# Copyright (c) 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "mount cgroups on /sys/fs/cgroup"
author          "chromium-os-dev@chromium.org"

start on starting boot-services
task

script
  mount -t tmpfs -o mode=755,noexec,nosuid,nodev none /sys/fs/cgroup
  mkdir -p /sys/fs/cgroup/cpu
  mount -t cgroup cgroup /sys/fs/cgroup/cpu -o cpu,noexec,nosuid,nodev
end script
