# init/clear-framebuffer.conf
# Copyright (c) 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description    "Clears the splash screen displayed by the boot-splash job"
author         "chromium-os-dev@chromium.org"

start on x-started

# When X starts, it copies the contents of the framebuffer to the root
# window.  We clear the framebuffer to make sure that it doesn't flash
# back onscreen when X exits later.
exec ply-image --clear 0x000000 --no-set-monitors
