# init/clean-shutdown.conf
# Copyright (c) 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Signal ready for clean shutdown to power manager"
author        "chromium-os-dev@chromium.org"

start on power-manager-clean-shutdown
task

script
  PM=org.chromium.PowerManager
  PM_PATH=/org/chromium/PowerManager
  dbus-send --type=signal --system $PM_PATH $PM.StartingShutdown
  CROMO=$(initctl status cromo)
  if [ "$CROMO" != "cromo stop/waiting" ]; then
    initctl stop cromo
  fi
  dbus-send --type=signal --system $PM_PATH $PM.CleanShutdown
end script
