# init/cleanup-disk.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "System software update service"
author "chromium-os-dev@chromium.org"

# when boot reaches login prompt, start up disk cleanup check.

start on stopped boot-complete

# Runs once and exits.  Cron also starts it every so often to check
# again.  The 'boot' argument tells the script where it was run from,
# and causes it to wait a while before doing anything.
exec /usr/sbin/chromeos-cleanup-disk boot
