# init/crash-reporter.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# chromeos-crash
#
# when boot reaches login prompt, starts chromeos crash handler.

start on stopped boot-complete

pre-start script
  # TODO(kmixter): Once crosbug.com/3030 is fixed, we can
  # remove this little section.
  STARTED_FILE=/tmp/crash-reporting-started
  [ -e $STARTED_FILE ] && exit 1
  touch $STARTED_FILE
end script

# This is run-once rather than a service.
task

exec /sbin/crash_reporter --init
