# init/crash-reporter.conf
# Copyright (c) 2011 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Initialize crash reporting services"
author          "chromium-os-dev@chromium.org"

# This job merely initializes its service and then terminates; the
# actual checking and reporting of crash dumps is triggered by an
# hourly cron job.
start on stopped boot-complete

pre-start script
  # TODO(kmixter): Once crosbug.com/3030 is fixed, we can
  # remove this little section.
  STARTED_FILE=/tmp/crash-reporting-started
  [ -e $STARTED_FILE ] && exit 1
  touch $STARTED_FILE
  mkdir -p /var/spool
end script

exec crash_reporter --init
