# init/cryptohomed.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# cryptohomed
#
# Starts the cryptohome daemon, which handles mounting and unmounting users'
# encrypted home directories.  Also supports offline login checks.
#

start on (started tcsd and started dbus)
stop on (starting reboot or starting halt)
respawn

script
export PKCS11_INIT_FLAG=
# TODO(crosbug.com/14277): Remove flag once this is the default and
# revert to expect fork/exec stanza for better upstart tracking of
# when cryptohome starts (in case a job wants to start on started
# cryptohome).
if [ -r "/home/chronos/.cryptohome-init-pkcs11" ]; then
  export PKCS11_INIT_FLAG="--cryptohome-init-pkcs11"
fi
exec /usr/sbin/cryptohomed --noclose ${PKCS11_INIT_FLAG}
end script
