# init/cryptohomed.conf
# Copyright (c) 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Start the cryptohome daemon"
author          "chromium-os-dev@chromium.org"

# Starts the cryptohome daemon, which handles mounting and unmounting users'
# encrypted home directories.  Also supports offline login checks.
# TODO(crosbug.com/32611): Add a proper 'stop on' stanza. Currently
# stays running to allow for unmounting during chromeos_shutdown.
start on started boot-services and started tcsd and started chapsd
respawn

# If attestation.epb still exists in its old location, move it to the new
# location where cryptohome will look for it.
pre-start script
  if [ -e /mnt/stateful_partition/home/.shadow/attestation.epb ]; then
    mv /mnt/stateful_partition/home/.shadow/attestation.epb \
      /mnt/stateful_partition/unencrypted/preserve/attestation.epb
  fi
end script

expect fork
exec cryptohomed --noclose

# Any long-running invocations of the client must first be terminated.
pre-stop exec pkill -TERM cryptohome
