# init/cryptohomed.conf
# Copyright (c) 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Start the cryptohome daemon"
author          "chromium-os-dev@chromium.org"

# Starts the cryptohome daemon, which handles mounting and unmounting users'
# encrypted home directories.  Also supports offline login checks.
start on started tcsd and started dbus and started chapsd
stop on starting reboot or starting halt
respawn

expect fork
exec cryptohomed --noclose

# Any long-running invocations of the client must first be terminated.
pre-stop exec pkill -TERM cryptohome
