# init/cryptohomed.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# cryptohomed
#
# Starts the cryptohome daemon, which handles mounting and unmounting users'
# encrypted home directories.  Also supports offline login checks.
#

start on (started tcsd and started dbus)
stop on (starting reboot or starting halt)
respawn

# TODO(fes): Remove the following line once Chrome calls TpmCanAttemptOwnership
post-start exec /bin/sh -c 'for i in $(seq 1 10); do sleep 1; \
  result=$(/usr/sbin/cryptohome --action=tpm_take_ownership); done;' &

expect fork
exec /usr/sbin/cryptohomed --noclose
