# init/date_proxy.conf
# Copyright (c) 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description  "Monitor the proxy settings used to access the date server"
author       "chromium-os-dev@chromium.org"

start on started system-services
stop on stopping system-services
respawn

script
  PROXY_FILE=/var/run/proxy
  dir=$(mktemp -d /tmp/date-proxy.$$.XXXXXXXXXXXX)
  trap "rm -rf $dir" EXIT
  fifo="$dir/fifo"
  mkfifo "$fifo"
  /sbin/minijail0 -u proxystate -g proxystate \
      /sbin/date-proxy-watcher 3>"$PROXY_FILE" 4>"$fifo"
  (while read x; do restart tlsdate; sleep 15; done) < "$fifo"
end script
