# init/dump-boot-stats.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# dump-boot-stats
#
# when login-prompt-ready is received, grabs the current uptime and
# disk stats.  Also sends uptime to metrics client for perf dashboard

start on login-prompt-ready

# This is run-once rather than a service.
task

script
  UPTIME_STATS=`cat /proc/uptime`
  DISK_STATS=`cat /sys/block/sda/stat`
  
  BOOT_TIME=`echo $UPTIME_STATS | cut -f 1 -d ' '`
  SECTORS_READ=`echo $DISK_STATS | cut -f 3 -d ' '`
  SECTORS_WRITTEN=`echo $DISK_STATS | cut -f 7 -d ' '`

  /usr/bin/metrics_client -t \
    Platform.BootTime $BOOT_TIME 1 30000 50 &
  /usr/bin/metrics_client \
    Platform.BootSectorsRead $SECTORS_READ 1 1000000 50 &
  /usr/bin/metrics_client \
    Platform.BootSectorsWritten $SECTORS_WRITTEN 1 10000 50 &
  
  echo $UPTIME_STATS > /tmp/uptime-login-prompt-ready
  echo $DISK_STATS > /tmp/disk-login-prompt-ready
end script
