# init/entd.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

start on started pkcs11
stop on stopping ui
respawn
normal exit 0
normal exit 1

env USER_ID=1000
env PKCS11_GID=208
env HOME=/home/chronos/user

# CHROMEOS_USER is exported by pkcs11-ready

pre-start script
  # If entd has trouble initializing the token, it will want to be restarted.
  # We run this script again so we can detect and fix a hosed tpm based token
  # before starting up again.
  /usr/sbin/fix_pkcs11_token.sh
end script

exec /sbin/minijail --uid=$USER_ID --gid=$PKCS11_GID -- \
     /usr/sbin/entdwife.sh -c start \
     --extension_path=${HOME}/Extensions \
     --utility="/etc/entd/policy-utils.js" \
     --username=${CHROMEOS_USER}
