# init/entd.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

start on start-user-session  # exports CHROMEOS_USER
stop on stopping ui
respawn
normal exit 0
normal exit 1

env HOME=/home/chronos/user

exec /usr/sbin/entd \
     --extension-path=${HOME}/Extensions \
     --utility="/etc/entd/policy-utils.js" \
     --pkcs11=glaptop \
     --username=${CHROMEOS_USER}
