# init/failsafe-delay.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Start the failsafe timer"
author          "chromium-os-dev@chromium.org"

# This job is used as the guarantee that the "failsafe" job will start.
# It's a separate job to allow failsafe to start on normal conditions
# before the timer here expires.

start on runlevel 2

exec sleep 30
