# init/failsafe.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "Job to trigger other jobs after system has been up for a while"
author      "chromium-os-dev@chromium.org"

# This job is used to trigger other jobs after the system has been up for
# an arbitrary amount of time, typically as backup in case some other job
# that's supposed to finish earlier doesn't run successfully.  For example:
#
#   start on stopped boot-complete or stopped failsafe
#
# can be used so that even if boot-complete doesn't run due to UI problems,
# services that enable the developer to SSH in or switch to a different VT
# so they can repair the system will still be executed eventually.

start on runlevel 2

task

exec sleep 30s
