# init/first.conf
# Copyright (c) 2011 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "First startup script"
author          "chromium-os-dev@chromium.org"

start on startup

script
  # Adjust OOM score for init to make everything unkillable
  # Chrome will adjust renderers to allow the system to kill them.
  echo -1000 > /proc/1/oom_score_adj
end script