# init/ibus.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Run ibus"
author        "chromium-os-dev@googlegroups.com"

start on login-prompt-ready
stop on stopping ui
respawn

env XAUTHORITY=/home/chronos/.Xauthority
env DISPLAY=:0.0
env MINIJAIL_FLAGS="--use-capabilities"
# might be overridden by --uid parameter in /sbin/session_manager_setup.sh
env USER_ID=1000
# ibus-daemon creates a temporary file in /home/chronos,
# not in /home/chronos/user.
env HOME=/home/chronos

script
# ibus checks this env variable to get the UID of the calling user.  It is
# meant for cases where userhelper is used, but it seems to be needed for us.
export USERHELPER_UID=${USER_ID}

exec /sbin/minijail ${MINIJAIL_FLAGS} --uid=${USER_ID} -- \
  /usr/bin/ibus-daemon --xim
end script

post-stop script

set +e
. /sbin/killers

# Ensure sub-processes of ibus-daemon are terminated. Note that these processes
# are not mananged by Upstart.
term_process "(candidate_w.*|ibus-[xem].*)"

# Since Chrome OS regenerates /var/lib/dbus/machine-id every time when the OS
# starts, ibus-daemon creates a new file in ${HOME}/.config/ibus/bus/ every time
# when it starts. It's better to remove them when ibus-daemon is terminated.
rm -f ${HOME}/.config/ibus/bus/*

end script
