# init/ibus.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Run ibus"
author        "chromium-os-dev@googlegroups.com"

start on login-prompt-ready
stop on stopping ui
respawn

env XAUTHORITY=/home/chronos/.Xauthority
env DISPLAY=:0.0
env MINIJAIL_FLAGS="--use-capabilities"
# might be overridden by --uid parameter in /sbin/session_manager_setup.sh
env USER_ID=1000

script
# ibus checks this env variable to get the UID of the calling user.  It is
# meant for cases where userhelper is used, but it seems to be needed for us.
export USERHELPER_UID=${USER_ID}

exec /sbin/minijail ${MINIJAIL_FLAGS} --uid=${USER_ID} -- \
  /usr/bin/ibus-daemon --xim
end script
