# init/iptables.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "Set iptables policies and add rules"
author      "chromium-os-dev@chromium.org"

# Make sure that we'll run eventually even if the UI doesn't come up correctly.
start on stopped boot-complete or stopped failsafe

# This is run-once rather than a service.
task

script

/sbin/iptables -P INPUT DROP
/sbin/iptables -P FORWARD DROP
/sbin/iptables -P OUTPUT DROP

# Accept everything on the loopback
/sbin/iptables -I INPUT -i lo -j ACCEPT
/sbin/iptables -I OUTPUT -o lo -j ACCEPT

# Accept return traffic inbound
/sbin/iptables -I INPUT  -m state --state ESTABLISHED,RELATED -j ACCEPT

# Accept icmp echo (NB: icmp echo ratelimiting is done by the kernel)
/sbin/iptables -A INPUT  -p icmp -j ACCEPT

# ssh
#for_test /sbin/iptables -A INPUT  -p tcp --dport 22 -j ACCEPT

# Accept new and return traffic outbound
/sbin/iptables -I OUTPUT  -m state --state NEW,ESTABLISHED,RELATED -j ACCEPT

end script
