# init/laptop-mode-resume.conf
# Copyright (c) 2013 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Update system power management settings after resume"
author        "chromium-os-dev@chromium.org"

start on system-resumed

exec laptop_mode auto
