# init/light-sensor.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Instantiate a light sensor device on I2C and load its driver"
author        "chromium-os-dev@chromium.org"

# Start this after udev since that appears to be the soonest we see
# an active I2C adapter available.

start on started udev
stop on starting halt or starting reboot

script

MACH=`uname -m`
if [ "$MACH" = "i686" -o "$MACH" = "x86_64" ]; then
  # Identify the I2C bus and address, and load the driver.
  echo tsl2563 0x29 > /sys/class/i2c-adapter/i2c-2/new_device
fi

end script
