# init/logout.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# logout
#
# Wrap the logout event
start on stopped ui
task

script
  # Indicate that we're no longer logged in
  rm -f /var/run/state/logged-in
  rm -f /var/run/state/windowmanager/initial-chrome-window-mapped
end script
