# init/logout.conf
# Copyright (c) 2011 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Mark that the user has logged out"
author          "chromium-os-dev@chromium.org"

# Use 'start on stopping' and 'task' to hold up restarting the
# session manager until we're done with marking the state change.
start on stopping ui
task

script
  # Indicate that we're no longer logged in
  rm -f /var/run/state/logged-in
  rm -f /var/run/state/windowmanager/initial-chrome-window-mapped

  # /var/run/headphone-jack-monitor.pid is created by
  #   /usr/sbin/headphone-jack-monitor.
  # When not be daemonized, no PID file will be created.
  if [ -f /var/run/headphone-jack-monitor.pid ]; then
     kill $(cat /var/run/headphone-jack-monitor.pid) && \
          rm /var/run/headphone-jack-monitor.pid;
  fi;

  # Pop user's network profile pushed in login.conf
  # NB: must quote ~chronos; don't want it expanded
  dbus-send --system --dest=org.chromium.flimflam --print-reply / \
    org.chromium.flimflam.Manager.PopProfile 'string:~chronos/flimflam'
end script
