# init/magic-sysrq.conf
# Copyright (c) 2011 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "Determine whether to allow the magic sysrq key to work"
author      "chromium-os-dev@chromium.org"

# Make sure that we'll run eventually even if the UI doesn't come up correctly.
start on stopped boot-complete or stopped failsafe

script
# If the dev switch is set or a shared user password has been set then
# we'll leave the magic sysrq enabled; else we'll disable it.
if ! (cros_boot_mode -s -m developer || grep -q '^chronos:[^*]' /etc/shadow)
then
  echo 0 > /proc/sys/kernel/sysrq
fi
end script
