# init/metrics_daemon.conf
# Copyright (c) 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Metrics collection daemon"
author          "chromium-os-dev@chromium.org"

# The metrics daemon is responsible for receiving and forwarding to
# chrome UMA statistics not produced by chrome.
#
# Starts after the crash reporter has finished so that kernel
# crashes can be appropriately detected and logged.
start on stopped crash-reporter
stop on starting halt or starting reboot
respawn

expect fork
exec metrics_daemon
