# init/metrics_daemon.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# The metrics daemon service used for passive metrics collection.

# Starts the metrics daemon after the crash reporter has finished so
# that kernel crashes can be appropriately detected and logged.
start on stopped crash-reporter
stop on starting halt or starting reboot

respawn
expect fork
exec /usr/bin/metrics_daemon
