# init/pkcs11.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Run pkcs#11 interface for the tpm on login"
author        "chromium-os-dev@googlegroups.com"

start on start-user-session
stop on stopping ui
# No respawning for now, but it is safe.
# respawn
expect fork

env USERNAME=chronos
env USER_ID=1000

env PKCS11_GID=208
env PKCS11_GROUP=pkcs11

export CHROMEOS_USER

# Hopefully this will pick up pkcsslotd but we may need to wait on it
pre-start script
  # Ensure the directories exist
  mkdir -p /var/lib/opencryptoki/tpm
  chown -R root:$PKCS11_GROUP /var/lib/opencryptoki/

  # Ensure that they point to the user volume
  [ -L /var/lib/opencryptoki/tpm/$USERNAME ] || \
    ln -sf /home/$USERNAME/user/.tpm /var/lib/opencryptoki/tpm/$USERNAME
  [ -L /var/lib/opencryptoki/tpm/root ] || \
    ln -sf ./$USERNAME /var/lib/opencryptoki/tpm/root

  # Always remove the old token entry.
  rm /var/lib/opencryptoki/pk_config_data || true

  # Creating this directory because if it's not there, token initialization
  # will neither create it nor populate it.
  mkdir -p /home/$USERNAME/user/.tpm/TOK_OBJ

  # Configure the tpm as a token
  pkcs_slot 0 tpm

  # Make sure the user can access their own data
  chown -R $USERNAME:$PKCS11_GROUP /home/$USERNAME/user/.tpm
end script

exec /sbin/minijail --uid=$USER_ID --gid=$PKCS11_GID -- /usr/sbin/pkcsslotd
