# init/pkcs11.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Run pkcs#11 interface for the tpm on login"
author        "chromium-os-dev@googlegroups.com"

start on start-user-session
stop on stopping ui
# No respawning for now, but it is safe.
# respawn
expect fork

env USER_ID=1000
env PKCS11_GID=208

export CHROMEOS_USER

pre-start script
  # We check the validity of the token before starting slotd, since the daemon
  # would be useless without it.  Note however that entd also does this, for
  # reasons described in entd.conf.
  /usr/sbin/fix_pkcs11_token.sh
end script

exec /sbin/minijail --uid=$USER_ID --gid=$PKCS11_GID -- /usr/sbin/pkcsslotd
