# init/pkcsslotd.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# PKCS#11 slot daemon, from opencryptoki.
# Brokers use of slots among applications, via a shm segment

start on starting tcsd
respawn

pre-start script
  rm -f /var/lib/opencryptoki/pk_config_data
  /usr/sbin/pkcs_slot 0 tpm
end script

expect fork
exec /usr/sbin/pkcsslotd
