# init/power.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# power
#
# When boot reaches login prompt and udev has started, set power
# related settings through laptop_mode.
# When gobi is reattached to the system (such as after resume),
# indicated by the signal gobi_modem, also rerun laptop_mode to ensure
# that the modem's usb power settings are correct.

start on (stopped boot-complete and started udev) or on gobi_modem

task

pre-start script
    if [ ! -d /var/run/laptop-mode-tools ]
    then
        mkdir -p /var/run/laptop-mode-tools
    fi
    touch /var/run/laptop-mode-tools/enabled
end script

exec /usr/sbin/laptop_mode force
