# init/power.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Initialize system power management settings"
author          "chromium-os-dev@chromium.org"

start on stopped boot-complete and started udev

pre-start script
  if [ ! -d /var/run/laptop-mode-tools ]
  then
    mkdir -p /var/run/laptop-mode-tools
  fi
  touch /var/run/laptop-mode-tools/enabled
end script

exec laptop_mode force
