# init/power.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Initialize system power management settings"
author          "chromium-os-dev@chromium.org"

start on starting system-services

pre-start script
  RUN_DIR=/var/run/laptop-mode-tools
  if [ ! -d ${RUN_DIR} ]; then
    mkdir -p ${RUN_DIR}
  fi
  touch ${RUN_DIR}/enabled
  if [ -f /home/chronos/.governor_* ]; then
    # Honor user-specified config by clone'ing the config dir
    # and applying the user config (for cpufreq); laptop_mode
    # uses the clone'd dir if present
    if [ ! -d ${RUN_DIR}/conf.d ]; then
      cp -r /etc/laptop-mode/conf.d ${RUN_DIR}
      C=${RUN_DIR}/conf.d/cpufreq.conf
      gov=$(basename /home/chronos/.governor_*)
      case "${gov}" in
      *_interactive)
        sed -i 's/_GOVERNOR=.*/_GOVERNOR=interactive/' ${C}
        # auto-enable input boost
        for s in BATT LM_AC NOLM_AC; do
          echo "${s}_CPU_INPUT_BOOST=1"
        done >> ${C}
        ;;
      *_ondemand)
        sed -i 's/_GOVERNOR=.*/_GOVERNOR=ondemand/' ${C}
        ;;
      *_performance)
        sed -i 's/_GOVERNOR=.*/_GOVERNOR=performance/' ${C}
        ;;
      *)
        logger -t "$UPSTART_JOB" "Unknown governor ${gov} requested"
        ;;
      esac
      cp /etc/laptop-mode/laptop-mode.conf ${RUN_DIR}/laptop-mode.conf
    fi
    logger -t "$UPSTART_JOB" "Setup customized laptop_mode config"
  fi
end script

exec laptop_mode force
