# init/powerm.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Start powerm after dbus has started"
author        "chromium-os-dev@chromium.org"

start on started dbus
stop on starting halt or starting reboot

env HOME=/home/chronos
env LOG_DIR=/var/log/power_manager
env USER_ID=1000

respawn
respawn limit 20 60  # if the job respawns 20 times in 60 seconds, stop trying.

script

mkdir -p ${LOG_DIR}
chown ${USER_ID} ${LOG_DIR}
cd ${HOME}
exec /usr/bin/powerm                                                      \
    --log_dir=${LOG_DIR}                                                  \
  > ${LOG_DIR}/powerm.out 2>&1; sleep 2

end script
