# init/power.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# power
#
# when boot reaches login prompt and udev has started, set power
# related settings

start on stopped boot-complete and started udev
# Respond to emit from cromo that gobi has reattached to system.
start on gobi_modem

#TODO: Add a stop on script
task

pre-start script
    if [ ! -d /var/run/laptop-mode-tools ]
    then
        mkdir -p /var/run/laptop-mode-tools
    fi
    touch /var/run/laptop-mode-tools/enabled
end script

exec /usr/sbin/laptop_mode force
