# init/pre-startup.conf
# Copyright 2014 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "System startup script"
author          "chromium-os-dev@chromium.org"

start on startup

script
  # The kernel and /sbin/init mount /proc, /sys and /dev before
  # we get here.
  mount -n -t tmpfs -o nodev,noexec,nosuid tmp /tmp

  mount -n -t tmpfs -o mode=0755,nodev,noexec,nosuid run /run
  mkdir -m 1777 /run/lock

  # Create static device nodes from modules.devname.
  static_node_tool || logger -t "$UPSTART_JOB" "static_node_tool failed."

  # With udev-208, static kernel modules are no longer automatically handled,
  # so load tun and ppp_async devices manually.
  modprobe tun || logger -t "$UPSTART_JOB" "modprobe tun failed."
  modprobe ppp_async || logger -t "$UPSTART_JOB" "modprobe ppp_async failed."
end script
