# init/pre-startup.conf
# Copyright 2014 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "System startup script"
author          "chromium-os-dev@chromium.org"

start on startup

script
  # The kernel and /sbin/init mount /proc, /sys and /dev before
  # we get here.
  mount -n -t tmpfs -o nodev,noexec,nosuid tmp /tmp
  # Mount debugfs as bootstat depends on /sys/kernel/debug
  mount -n -t debugfs -o nodev,noexec,nosuid,mode=0750,uid=0,gid=debugfs-access \
    debugfs /sys/kernel/debug

  mount -n -t tmpfs -o mode=0755,nodev,noexec,nosuid run /run
  mkdir -m 1777 /run/lock

  # bootstat writes timings to both tmpfs and debugfs.
  bootstat pre-startup
end script
