# init/procps.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# procps
#
# when boot reaches login prompt, runs procps which sets sysctls.

start on stopped boot-complete

# This is run-once rather than a service.
task

script
  # TODO: Move all sysctl settings into one file.
  for i in /etc/sysctl.conf /etc/sysctl.d/* ; do
    /sbin/sysctl -q -p "$i"
  done
end script
