# init/pulseaudio.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Run pulseaudio when the system starts"
author        "chromium-os-dev@googlegroups.com"

# pulseaudio needs udev to be up and running
start on started udev
stop on starting halt
respawn

env USERNAME=chronos
env XAUTHORITY=/home/chronos/.Xauthority
env DISPLAY=:0.0
# TODO(ajwong): pulseaudio wants CAP_NICE.  Readd --use-capabilities after
# resolving that dependency.
env MINIJAIL_FLAGS=""
env PULSE_STATE_PATH=/home/chronos/.pulse
env PULSE_RUNTIME_PATH=/var/run/pulse
env USER_ID=1000

pre-start script
  mkdir -p "${PULSE_RUNTIME_PATH}"
  chown -R chronos:chronos "${PULSE_RUNTIME_PATH}"
end script

expect fork
exec /sbin/minijail ${MINIJAIL_FLAGS} --uid=${USER_ID} -- \
  /usr/bin/pulseaudio --log-target=syslog --start

# Clean up our run directory
post-stop script
  rm -rf "${PULSE_RUNTIME_PATH}"
end script
