# init/recovery.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# description "Recovery service"
# author "chromium-os-dev@googlegroups.com"

start on started udev

# TODO(tgao): replace debug output and auto-shutdown with a friendly UI
script
    if [ -e /mnt/stateful_partition/.recovery ]; then
      # use "sh -x" to output extensive debug output to console for now
      sh -x /usr/sbin/chromeos-install --recovery --run_as_root --yes

      # shut down and power off the device (with 1 minute delay)
      shutdown -h +1
    else
      echo "Regular boot. Skip recovery" >> /var/log/recovery.log
    fi
end script
console output
