# init/screen-locker.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Run the screen locker for the given logged in user"
author        "chromium-os-dev@chromium.org"

# start-user-session is only emitted from the new session_manager code.
# chromeos-xsession (the old way of starting the WM) is not run on that
# code pathway.
start on start-user-session
stop on stopping ui

# sadly, these can't reference each other.
env DATA_DIR=/home/chronos
env HOME=/home/chronos/user
env XAUTHORITY=/home/chronos/.Xauthority
env DISPLAY=:0.0

# might be overridden by --uid parameter in /sbin/session_manager_setup.sh
env USER_ID=1000

respawn

script
# May have been created by window manager.
mkdir -p "$HOME/log"

# Define some max.  Don't need that many for debugging.
local max_logs=5

# All logs are tagged with the date.
local log_tag=`date +%Y%m%d-%H%M%S`
local log_count=`ls -t "$HOME/log/screen-locker"* | wc -l`

# Don't allow logs to grow past max log count.
if [ $log_count -ge $max_logs ] ; then
  local oldest_log=`ls -rt "$HOME/log/screen-locker"* | head -1`
  rm "$oldest_log"
fi

exec /sbin/minijail --uid=${USER_ID} -- \
  /usr/bin/xscreensaver -no-splash -v -log "$HOME/log/screen-locker.$log_tag"

end script
