# init/sleep.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# sleep
#
# used for delayed start

start on runlevel 2

task

exec sleep 30s
