# init/syslog.conf
# Copyright (c) 2011 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.
#

description   "Syslog daemon"
author        "chromium-os-dev@chromium.org"

start on starting boot-services
stop on stopping boot-services

pre-start script
  SYSLOG_RUN="/var/run/rsyslog"
  KMSG_PIPE="${SYSLOG_RUN}/kmsg"

  if [ ! -d "${SYSLOG_RUN}" ]; then
    install --owner=syslog --group=syslog -d "${SYSLOG_RUN}"
  fi
  if [ ! -p "${KMSG_PIPE}" ]; then
    mkfifo -m 600 "${KMSG_PIPE}"
    chown syslog:syslog "${KMSG_PIPE}"
  fi
end script

respawn
expect fork
exec rsyslogd -c4 -f /etc/rsyslog.chromeos
