# init/syslog.conf
# Copyright (c) 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Syslog daemon"
author        "chromium-os-dev@chromium.org"

start on starting boot-services
stop on stopping boot-services
respawn

pre-start script
  touch /var/log/messages
  chown syslog /var/log/messages
end script

exec rsyslogd -c4 -n -f /etc/rsyslog.chromeos
