# init/tcsd.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# Trousers daemon, which talks to the TPM (or the TPM emulator).

# Obviously, this start rule works only if we start the emulator.
# TODO(semenzato): Nathan Williams's suggestion: we create an event name as the
# API ("tpm-device-added" or something), and arrange the udev rules for both
# the real device and the emulator to emit the same event.  We can use
# initctl-emit.
start on started tpm-emulator
stop on starting halt or starting reboot

respawn
expect fork
exec /usr/sbin/tcsd
