# init/tpm-emulator.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# TPM emulator.

start on runlevel 2
stop on starting halt or starting reboot

pre-start script
  # The tpmd_dev udev rule creates the tpm devices when this module
  # is inserted.  Since there is no actual hardware I don't know how
  # else to get this started (suggestions welcomed to semenzato@).
  modprobe tpmd_dev
end script

respawn
expect fork
exec /usr/sbin/tpmd
