# init/tpm-probe.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# Detects TPM hardware, and starts the emulator (if needed) and the trousers
# daemon.

start on started syslog

script
  # If there is a TPM, and the firmware detected it, and the module is compiled
  # in the kernel, or has been already loaded in other ways, then the module
  # has created a sysfs entry, and we're done.
  #
  # Otherwise there are two cases: there is a TPM but it hasn't been discovered
  # yet (possibly because we're using a TPM-agnostic firmware for testing), so
  # we force its detection; or there is no TPM, and we start the emulator.

  if ! test -e /sys/devices/platform/tpm_tis && \
     ! modprobe tpm_tis force=1 interrupts=0; then
    rmmod tpm tpm_bios # added by modprobe, but not needed by the emulator
    start tpm-emulator
  fi

  # At this point we either have the hardware or the emulator, so start tcsd.
  start tcsd
end script
