# init/udev.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# udev
#
# when dump-boot-stats has run, starts udev.

start on stopped dump-boot-stats

respawn
expect fork
exec /sbin/udevd --daemon

pre-start script
  # udev uses inotify to monitor the dynamic rules directory
  # /dev/.udev/rules.d, however this directory must exist before udevd
  # runs, or the inotify won't find the directory even if it is created
  # later.  This is an issue with udev-146.  Udev versions 156 and later
  # do not exhibit this problem, and for those this script is not necessary.
  if [ ! -d /dev/.udev/rules.d ]; then
    mkdir -p /dev/.udev/rules.d
  fi
end script

post-start script
  # In some systems hotplug is used before udev starts.  We need to disable it
  # so it doesn't fight with udev.
  if [ -e /proc/sys/kernel/hotplug ]; then
    echo "" > /proc/sys/kernel/hotplug
  fi

  /sbin/udevadm trigger
  /sbin/udevadm settle
end script
