# init/udev-addon.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# udev-addon
#
# Kick udev to find everything else either when boot reaches login prompt
# or after a delay period specified in the sleep job.  Starting after the
# sleep job ensures that there is networking even if UI fails.

start on stopped boot-complete or stopped sleep

task

script
  /sbin/udevadm trigger
  /sbin/udevadm settle
end script

post-start script
  wlanfile="/var/lib/chromeos-aliases.conf"
  # Create/Update alias for the built-in WIFI.
  # wlan0 alias is used to preload WIFI driver in udev.conf
  # TODO(ggg): flimflam knows the network device name. It should tell us.
  d=$(ethtool -i wlan0 | awk '/driver:/ { print $2 }')
  if [ ! -z "$d" ] ; then
    # try to determine previous alias
    wlan=""
    if [ -s $wlanfile ] ; then
      wlan=$(awk '/alias wlan0/ { print $3 ; exit 0 }' $wlanfile)
      if [ ! -z "$wlan" ] ; then
        sed -i '/wlan0/d' $wlanfile
      fi
    fi
    if modinfo $d 1>/dev/null 2>&1 ; then
      # $d is a module.
      # if no alias existed or existing one was different, write alias out
      if [ ! -z "$wlan" -o "$wlan" != "$d" ] ; then
        echo "alias wlan0 $d" > $wlanfile
      fi
    else
      # $d is not a module
      rm -f $wlanfile
      /usr/bin/logger "udev-addon:  $wlanfile removed. $d is not a module."
    fi
  fi
end script
