# init/udev-trigger-input.conf
# Copyright (c) 2011 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Replay udev input events"
author          "chromium-os-dev@chromium.org"

# Replay input subsystem kernel events that arrived before udevd was started.
# These events are needed so X can initialize its input devices on startup.
# Other udev events will get replayed by udev-trigger.conf
start on starting boot-services

script
  udevadm trigger --subsystem-match="input"
  # settle required only for powerm
  exec udevadm settle
end script
