# init/udev.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# udev
#
# Start udev early so flimflam can be notified of "wlan0"
# The general "trigger" for remaining devices is now in udev-addon.conf

start on runlevel 2

respawn
expect fork
exec /sbin/udevd --daemon

pre-start script
  # udev uses inotify to monitor the dynamic rules directory
  # /dev/.udev/rules.d, however this directory must exist before udevd
  # runs, or the inotify won't find the directory even if it is created
  # later.  This is an issue with udev-146.  Udev versions 156 and later
  # do not exhibit this problem, and for those this script is not necessary.
  if [ ! -d /dev/.udev/rules.d ]; then
    mkdir -p /dev/.udev/rules.d
  fi 

  # If driver for built-in WIFI is a module, udev-addon will record that.
  # modprobe gets very unhappy if it can't find an alias.
  d="wlan0"
  f="/var/lib/chromeos-aliases.conf"
  if [ -s "$f" ] ; then
    if ! modprobe -C $f $d 1>/dev/null 2>&1 ; then
        # modprobe didn't like the conf file for some reason.
        # move it out of the way so udev-addon.conf will make a new one.
        mv $f $f.bad
	/usr/bin/logger "udev pre-start: modprobe using $f failed.  Renamed with .bad suffix."
    fi
  fi
end script
