# init/udev.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Start the udev daemon"
author          "chromium-os-dev@chromium.org"

# Start udev early so flimflam can be notified of "wlan0".  The
# general "trigger" for remaining devices is now in the udev-addon
# job
start on runlevel 2

respawn
expect fork

pre-start script
  # udev uses inotify to monitor the dynamic rules directory
  # /dev/.udev/rules.d; however, this directory must exist before
  # udevd runs, or the inotify won't find the directory even if it
  # is created later.  This is an issue with udev-146.  Udev
  # versions 156 and later do not exhibit this problem, and for
  # those this test is not necessary.
  if [ ! -d /dev/.udev/rules.d ]; then
    mkdir -p /dev/.udev/rules.d
  fi

  # If the driver for built-in WIFI is a module, the udev-addon job
  # will have recorded it in an aliases file on the previous boot
  # cycle.  Try and start that same device now, so that flimflam can
  # get started bringing up the wireless LAN sooner.
  #
  # Note that although a change is extremely unlikely, we don't
  # consider failure to probe the module here fatal.
  dev_alias=wlan0
  aliasfile=/var/lib/chromeos-aliases.conf
  if [ -s $aliasfile ] ; then
    if ! modprobe -C $aliasfile $dev_alias >/dev/null 2>&1 ; then
      # If modprobe didn't like the alias file for some reason, just
      # move it out of the way.  The udev-addon job will make a new
      # one.
      mv $aliasfile $aliasfile.bad
      MESSAGE="udev pre-start: modprobe using $aliasfile failed."
      MESSAGE="$MESSAGE  Renamed with .bad suffix."
      logger "$MESSAGE"
    fi
  fi
end script

exec udevd --daemon
