# init/udev-addon.conf
# Copyright (c) 2011 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Trigger udev events"
author          "chromium-os-dev@chromium.org"

# We must run eventually even if the UI doesn't come up correctly,
# because network configuration depends on this job.
start on starting failsafe

script
  udevadm trigger
  udevadm settle

  wlanfile=/var/lib/chromeos-aliases.conf
  # Create/Update alias for the built-in WIFI.
  # wlan0 alias is used to preload WIFI driver in udev.conf
  # TODO(ggg): flimflam knows the network device name. It should tell us.
  dev=$(ethtool -i wlan0 | awk '/driver:/ { print $2 }')
  if [ -n "$dev" ] ; then
    # First, clear out any previous alias, in case there's been a
    # change.
    if [ -s $wlanfile ] ; then
      sed -i '/alias wlan0/d' $wlanfile
    fi
    # Now, if we can create a valid alias from $dev, do so.
    if modinfo $dev >/dev/null 2>&1 ; then
      # $dev is a module.
      echo "alias wlan0 $dev" >>$wlanfile
    else
      # $dev is not a module
      logger "udev-addon:  no alias created for wlan0. $dev is not a module."
    fi
  fi
end script
