# init/udev-addon.conf
# Copyright (c) 2011 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Trigger udev events"
author          "chromium-os-dev@chromium.org"

# We must run eventually even if the UI doesn't come up correctly,
# because network configuration depends on this job.
start on stopped boot-complete or stopped failsafe

script
  udevadm trigger
  udevadm settle
end script

post-start script
  wlanfile="/var/lib/chromeos-aliases.conf"
  # Create/Update alias for the built-in WIFI.
  # wlan0 alias is used to preload WIFI driver in udev.conf
  # TODO(ggg): flimflam knows the network device name. It should tell us.
  d=$(ethtool -i wlan0 | awk '/driver:/ { print $2 }')
  if [ ! -z "$d" ] ; then
    # try to determine previous alias
    wlan=""
    if [ -s $wlanfile ] ; then
      wlan=$(awk '/alias wlan0/ { print $3 ; exit 0 }' $wlanfile)
      if [ ! -z "$wlan" ] ; then
        sed -i '/wlan0/d' $wlanfile
      fi
    fi
    if modinfo $d 1>/dev/null 2>&1 ; then
      # $d is a module.
      # if no alias existed or existing one was different, write alias out
      if [ ! -z "$wlan" -o "$wlan" != "$d" ] ; then
        echo "alias wlan0 $d" > $wlanfile
      fi
    else
      # $d is not a module
      rm -f $wlanfile
      logger "udev-addon:  $wlanfile removed. $d is not a module."
    fi
  fi
end script
