# init/update-engine.conf
# Copyright (c) 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "System software update service"
author          "chromium-os-dev@chromium.org"

start on starting system-services
stop on stopping system-services
respawn

expect fork

# Runs the daemon at low/idle IO priority so that updates don't
# impact system responsiveness.
exec ionice -c3 update_engine
