# init/update-engine.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "System software update service"
author "chromium-os-dev@chromium.org"
#
# When boot reaches login prompt and the D-Bus daemon is ready, start up
# software update check.

start on stopped boot-complete and started dbus
stop on starting halt or starting reboot

respawn
expect fork

# Runs the daemon at low/idle IO priority so that updates don't impact
# system responsiveness.
exec /usr/bin/ionice -c3 /usr/sbin/update_engine
