# init/update-engine.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# description "System software update service"
# author "chromium-os-dev@googlegroups.com"
#
# when dump-boot-stats has run, start up software update check.

start on stopped dump-boot-stats
stop on starting halt or starting reboot

respawn
expect fork
exec /usr/sbin/update_engine
