# init/ureadahead.conf
# Copyright (c) 2013 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Read files required for boot in advance"
author          "chromium-os-dev@chromium.org"

# Start ureadahead as early as possible. Note that it depends on some of the
# operations in chromeos_startup. Notably, it requires /var to be mounted,
# because the pack file is located in /var/lib/ureadahead, and ureadahead
# tracing requires /sys/kernel/debug be mounted.
start on stopping startup
stop on starting boot-complete

# Give ureadahead extra time after SIGTERM to write out the pack file.
kill timeout 30

pre-start exec mkdir -p -m 0700 /var/lib/ureadahead

exec ureadahead
