# init/vpd-log.conf
# Copyright (c) 2011 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Dumps VPD information for chrome://system"
author          "chromium-os-dev@chromium.org"

# VPD (Vital Product Data) information is shown on chrome://system for
# user, feedback, and partner's customer service center.
#
# This job will read VPD from BIOS and log it into /var/log/vpd_2.0.txt.
# It should be executed at least once some time after boot and before the
# user types "chrome://system" in the chrome location bar.
#
# This job depends on system-services so as to keep to it off of the
# boot critical path, because it may invoke the flashrom utility to
# access BIOS, which is slow.
start on starting system-services

script
  if [ -x /usr/sbin/dump_vpd_log ]; then
    /usr/sbin/dump_vpd_log
  fi
end script
