# init/vt-switching.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "Determine whether to permit switching between virtual terminals"
author      "chromium-os-dev@googlegroups.com"

start on stopped boot-complete

script
STATE_DIR='/var/run/state'
ALLOW_VT_SWITCHING_FILE='allow-vt-switching'

ALLOW_VT_SWITCHING_PATH="${STATE_DIR}/${ALLOW_VT_SWITCHING_FILE}"

mkdir -p "$STATE_DIR"

# If the dev switch is set or a shared user password has been set, then
# touch a file to let X know that it should allow switching.
if cros_boot_mode -s -m developer || grep -q '^chronos:[^*]' /etc/shadow; then
  touch "$ALLOW_VT_SWITCHING_PATH"
else
  rm -f "$ALLOW_VT_SWITCHING_PATH"
fi
end script
