# init/windowmanager_not_logged_in.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Run the window manager before the user has logged in"
author        "chromium-os-dev@chromium.org"

# The window manager gets run twice. Once before the user logs in (this script)
# and again after the user logs in (windowmanager.conf). This is done so that
# the window manager knows to disable certain operations we don't want to allow
# when the user isn't logged in (such as creating a terminal).
start on x-started
stop on start-user-session or stopping ui

# NOTE: /home/chronos/user has not yet been mounted.
env HOME=/home/chronos
env LOG_DIR=/var/log/window_manager
env XAUTHORITY=/home/chronos/.Xauthority
env DISPLAY=:0.0

# might be overridden by --uid parameter in /sbin/session_manager_setup.sh
env USER_ID=1000

env PULSE_RUNTIME_PATH=/var/run/pulse

respawn

script

WM="/usr/bin/chromeos-wm"
IMAGES="/usr/share/chromeos-assets/images"

mkdir -p ${LOG_DIR}
chown ${USER_ID} ${LOG_DIR}

cd ${HOME}
exec /sbin/minijail --uid=${USER_ID} --                                   \
  "${WM}"                                                                 \
    --hotkey_overlay_image_dir="${IMAGES}"                                \
    --log_dir="${LOG_DIR}"                                                \
    --panel_anchor_image="${IMAGES}/panel_anchor.png"                     \
    --panel_dock_background_image="${IMAGES}/panel_dock_bg.png"           \
    --shadow_image_dir="${IMAGES}"                                        \
    --logged_in=false
end script
