# init/wpasupplicant.conf
# Copyright (c) 2009 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# wpa_supplicant
#
# starts wpa_supplicant once dbus has been started.

start on started dbus

respawn

# make sure the control interface dir exists and has the right owner
pre-start script
  mkdir -p -m 0750 /var/run/wpa_supplicant
  chown wpa:wpa /var/run/wpa_supplicant
end script

# use minijail (drop root, keep CAP_NET_ADMIN|CAP_NET_RAW)
# enable dbus, syslog msgs, wpa_cli over unix domain socket
exec /sbin/minijail0 -u wpa -g wpa -c 3000 -- \
  /usr/sbin/wpa_supplicant -u -s -O/var/run/wpa_supplicant
