# libwebserv/webservd/etc/init/webservd.conf
# Copyright 2015 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Brillo WebServer Daemon"
author          "chromium-os-dev@chromium.org"

start on starting system-services
stop on stopping system-services
respawn

pre-start script
  mkdir -m 0755 -p /var/log/webservd
  chown webservd:webservd /var/log/webservd
end script

env WEBSERVD_LOG_LEVEL=0

# Minijail actually forks off our desired process.
expect fork

exec minijail0 -i -g webservd -u webservd /usr/bin/webservd \
    --v="${WEBSERVD_LOG_LEVEL}"
