# cros-disks/cros-disks.conf
# Copyright (c) 2011 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# cros-disks upstart job

env CROS_DISKS_LOG_DIR=/var/log/cros-disks
env CROS_DISKS_UID=213
env CROS_DISKS_GID=213

start on started dbus
stop on starting halt or starting reboot

respawn
expect fork

pre-start script
  mkdir -p -m 0755 "${CROS_DISKS_LOG_DIR}"
  chown -R cros-disks:cros-disks "${CROS_DISKS_LOG_DIR}"
end script

exec /sbin/minijail --uid="${CROS_DISKS_UID}" --gid="${CROS_DISKS_GID}" -- /opt/google/cros-disks/disks 

