# imageloader/init/imageloader-shutdown.conf
# Copyright 2017 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "cleanup of flash component updater"
author "chromium-os-dev@chromium.org"

start on stopped ui
task

# Cleanup all imageloader mounts so other filesystems can unmount cleanly.
exec /usr/sbin/imageloader --unmount_all
