# bluetooth/init/upstart/bluetooth-setup.conf
# Copyright 2018 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Prepare LE kernel splitter for Bluetooth daemons"
author        "chromium-os-dev@chromium.org"

start on started system-services
task

exec /usr/bin/bluetooth-setup.sh
