# arc/vm/init/arcvm-server-proxy.conf
# Copyright 2018 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Start the server proxy for arcvm"
author        "chromium-os-dev@chromium.org"

# TODO(yusukes): Stop depending on arcvm.conf.

# arcvm-server-proxy should outlive the arcvm conceptually so that
# clientproxy in the guest can connect to the server-proxy.
start on starting arcvm or starting vm_concierge
stop on stopped arcvm or stopped vm_concierge

pre-start script
  logger -t "${UPSTART_JOB}" "Pre-start arcvm-server-proxy"
  grep -qs  "1" /run/chrome/is_arcvm || exit 1
  # Make sure the vsock module is loaded.
  exec modprobe -q vhost-vsock
end script

# Use minimalistic-mountns profile.
# -e for a new network namespace.
# -p -I for a new PID namespace and run the process as init (no need to fork).
# -l for a new IPC namespace.
# --uts for UTS namespace to isolate from host / domain names.
# -N for freeze cgroup settings.
# /mnt is the mount point of the fuse file system.
# /run/chrome/arc_bridge.sock is the socket connected to the ArcBridgeService
# in Chrome browser process.
exec minijail0 \
  --profile=minimalistic-mountns \
  -e \
  -p -I \
  -l \
  --uts \
  -N \
  -k "tmpfs,/mnt,tmpfs,MS_NOSUID|MS_NODEV|MS_NOEXEC" \
  -k "tmpfs,/run,tmpfs,MS_NOSUID|MS_NODEV|MS_NOEXEC" \
  -b /dev/fuse \
  -b /run/chrome/arc_bridge.sock \
  -- /usr/bin/arcvm_server_proxy /mnt

post-stop exec logger -t "${UPSTART_JOB}" "Post-stop arcvm-server-proxy"
