# camera/hal_adapter/init/camera-halv3-adapter.conf
# Copyright 2017 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Start camera HAL v3 adapter service"
author        "chromium-os-dev@chromium.org"

start on starting system-services
stop on stopping system-services
expect fork
respawn
respawn limit 10 5

env SECCOMP_POLICY_FILE=/usr/share/policy/camera-halv3-adapter.policy

# The user arc-camera and the binary arc_camera3_service have arc in its
# naming for historical reasons (carried from camera HAL v1). For camera
# HAL v3 they are not just ARC-specific: Chrome hosts the CameraHalDispatcher
# Mojo proxy; arc_camera3_service forks a child process to host the camera
# HAL adapter and registers to CameraHalDispatcher Mojo proxy. Camera clients
# can then ask CameraHalDispatcher Mojo proxy to get the handle to connect
# to the camera HAL adapter.

# Need -G to be in video group to access /dev/video*
exec minijail0 -i -u arc-camera -g arc-camera -G \
    -e -N -p -r -v -l \
    -n -L -S "${SECCOMP_POLICY_FILE}" \
    -- /usr/bin/arc_camera3_service
