# dlcservice/dlcservice.conf
# Copyright 2018 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Downloadable content service daemon"
author          "chromium-os-dev@chromium.org"

start on starting system-services
stop on stopping system-services

exec minijail0 -u dlcservice -g dlcservice -G -n -S \
  /usr/share/policy/dlcservice-seccomp.policy /usr/sbin/dlcservice
