# dlcservice/dlcservice.conf
# Copyright 2018 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Downloadable content service daemon"
author          "chromium-os-dev@chromium.org"

start on starting system-services
stop on stopping system-services
respawn
respawn limit 3 10  # if the job respawns 3 times in 10 seconds, stop trying.

pre-start script
  # Initializes DLC module image root directory.
  dlc_image_root_path="/var/cache/dlc"
  mkdir -p "${dlc_image_root_path}"
  chmod 0755 "${dlc_image_root_path}"
  chown dlcservice:dlcservice "${dlc_image_root_path}"
end script

exec minijail0 -u dlcservice -g dlcservice -G -n -S \
  /usr/share/policy/dlcservice-seccomp.policy /usr/sbin/dlcservice
