# midis/init/midis.conf
# Copyright 2017 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Chrome OS MIDI server"
author          "chromium-os-dev@chromium.org"

start on start-arc-instance-for-login-screen or start-arc-instance
stop on stop-arc-instance
respawn

# Allow the MIDI server real time priority.
limit rtprio 20 20

exec minijail0 -l -p -r -v -N -n \
  -S /usr/share/policy/midis-seccomp.policy \
  -c 0 -u midis -g midis -G -- /usr/bin/midis
