# midis/init/midis.conf
# Copyright 2017 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Chrome OS MIDI server"
author          "chromium-os-dev@chromium.org"

start on starting system-services
stop on stopping system-services
respawn
expect fork

# Allow the MIDI server real time priority.
limit rtprio 20 20

exec minijail0 -l -p -r -v -N \
  -c 0x200000 -u midis -g midis -G -- /usr/bin/midis
