# cecservice/share/cecservice.conf
# Copyright 2018 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "cecservice"
author        "chromium-os-dev@chromium.org"

start on started system-services
stop on stopping system-services
respawn
respawn limit 3 10  # if the job respawns 3 times in 10 seconds, stop trying.

expect fork

# -u run as user cecservice
# -g run as group cecservice
# -i so Upstart tracks the daemon proces
# -l enter new IPC namespace
# -N enter new cgroup namespace
# -p enter new pid namespace
# -n set no_new_privs
# -w create and join a new anonymous session keyring
# -C run in a chroot
# -t mount a new /tmp inside chroot
# -b directories to mount into chroot
#
exec minijail0 -u cecservice -g cecservice -i -l -N -p -n -w -t \
     -C /run/cecservice_chroot -b /dev,/dev -b /sys,/sys \
     -b /usr/sbin,/usr/sbin -b /run/udev,/run/udev \
     -b /run/dbus,/run/dbus -b /tmp,/tmp --uts \
     -- /usr/sbin/cecservice
