# shill/init/shill-start-user-session.conf
# Copyright 2018 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Perform connection-manager-specific session setup tasks."
author          "chromium-os-dev@chromium.org"

# We're a task to ensure that the connection manager is ready to continue
# before session manager calls this step done.
start on start-user-session
task

# CHROMEOS_USER is an upstart-provided variable setting passed from
# session_manager via the 'start-user-session' event.
import CHROMEOS_USER
exec shill_login_user ${CHROMEOS_USER}

