# authpolicy/etc/init/authpolicyd.conf
# Copyright 2016 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Authpolicy daemon"
author          "chromium-os-dev@chromium.org"

start on starting ui
stop on stopping ui
respawn

# Minijail actually forks off our desired process.
expect fork

exec minijail0 -i -I -p -e -l -r -v -u authpolicyd -g authpolicyd /usr/sbin/authpolicyd

# Wait for daemon to claim its D-Bus name before transitioning to started.
post-start exec gdbus wait --system --timeout 15 org.chromium.authpolicy
