# authpolicy/etc/init/authpolicyd.conf
# Copyright 2016 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Authpolicy daemon"
author          "chromium-os-dev@chromium.org"

start on starting ui
stop on stopping ui
respawn

# krb5_machine.keytab and smb.conf gonna live there. All the other files will go
# into /tmp new mounted on each start.
pre-start script
  mkdir -m 0700 -p /var/lib/authpolicyd
  chown -R authpolicyd:authpolicyd /var/lib/authpolicyd
end script


# Minijail actually forks off our desired process.
expect fork

# -e is not specified cause the service needs to connect to AD server to join
# a domain, authenticate users, fetch user and device policies.
exec minijail0 -i -I -p -l -r -v -t -u authpolicyd -g authpolicyd /usr/sbin/authpolicyd

# Wait for daemon to claim its D-Bus name before transitioning to started.
post-start exec gdbus wait --system --timeout 15 org.chromium.AuthPolicy
