# container_utils/broker-service.conf
# Copyright 2016 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "broker_service (adb)"

start on socket PROTO=unix SOCKET_PATH=/run/broker_service/adb

# broker_service needs to be chronos due to permission_broker's DBus
# policy.
exec minijail0 -u chronos -g chronos -G /usr/bin/broker_service
