# container_utils/broker-service-pre-upstart-socket-bridge.conf
# Copyright 2016 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "Prepare upstart-socket-bridge for broker_service"

start on starting upstart-socket-bridge

env BROKER_SERVICE_SOCKET_DIR=/run/broker_service

script
  mkdir -p -m 0755 "${BROKER_SERVICE_SOCKET_DIR}"
end script
