# power_manager/init/report-power-metrics.conf
# Copyright (c) 2014 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Report resume metrics once per boot"
author        "chromium-os-dev@chromium.org"

start on started system-services

# Report resume metrics after boot in case a shutdown occurred while
# the computer was suspended
exec send_metrics_on_resume -b
