# userfeedback/init/firmware-version.conf
# Copyright (c) 2011 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Launch firmware version scripts to populate chrome://system"
author        "chromium-os-dev@chromium.org"

start on started system-services
script
  /usr/share/userfeedback/scripts/bios_info
  /usr/share/userfeedback/scripts/ec_info
end script
