# userfeedback/init/memory-spd-info.conf
# Copyright (c) 2011 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Launch memory SPD info scripts to populate chrome://system"
author        "chromium-os-dev@chromium.org"

start on started system-services
script
  /usr/share/userfeedback/scripts/memory_spd_info
end script
