# userfeedback/init/mosys-info.conf
# Copyright (c) 2011 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Launch mosys info dump script to populate chrome://system"
author        "chromium-os-dev@chromium.org"

start on started system-services
exec /usr/share/userfeedback/scripts/mosys_info
