# crosdns/init/crosdns.conf
# Copyright 2018 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description      "CrosDns daemon for /etc/hosts modifications"
author           "chromium-os-dev@chromium.org"

# Starts the crosdns daemon which provides a service for making
# modifications to the /etc/hosts file.
start on starting vm_concierge
stop on stopped vm_concierge
respawn

pre-start script
  # Create the directory we use for storing our copy of /etc/hosts and then
  # copy the existing version there as a baseline and bind mount it on top of
  # the existing one.
  mkdir -p -m 0755 /run/crosdns/
  chown crosdns:crosdns /run/crosdns/
  cp /etc/hosts /run/crosdns/hosts
  chown crosdns:crosdns /run/crosdns/hosts
  chmod 0644 /run/crosdns/hosts
  mount -o bind /run/crosdns/ /etc/hostsdir
end script

post-stop script
  umount --lazy /etc/hostsdir
end script

# use minijail (no caps, ipc namespace, pid namespace,
# mount namespace, remount proc, tmpfs, UTS namespace, drop root,
# set no_new_privs, set seccomp filter)
script
# TODO(jkardatzke): Add the seccomp policy and minijail.
#  exec minijail0 -c 0 -l -p -v -r -t --uts -u crosdns -g crosdns \
#      -n -S /usr/share/policy/crosdns-seccomp.policy \
#      /usr/sbin/crosdns
   exec minijail0 -u crosdns -g crosdns /usr/sbin/crosdns
end script
