# arc/setup/etc/arc-boot-continue.conf
# Copyright 2017 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Turn the container into a fully functional one"
author        "chromium-os-dev@chromium.org"

start on continue-arc-boot
stop on stop-arc-instance or stopping ui

env LOGFILE=/var/log/arc-boot-continue.log
env ENVFILE=/etc/init/arc-setup-env

# Do everything in pre-start to block session_manager. See arc-setup.conf
# for more details.
pre-start script
  bootstat android-start
  exec > $LOGFILE 2>&1
  echo "$(date --rfc-3339=ns): Pre-start arc-boot-continue"
  set -x
  . $ENVFILE
  /usr/sbin/arc-setup --boot-continue
  initctl start -n arc-sdcard
  exec bootstat arc-setup-end
end script

# This file should have neither start nor post-start scripts. See arc-setup.conf
# for more details.

# This file doesn't have a post-stop script. arc-setup-for-login-screen.conf's
# cleans things up.
