# arc/setup/etc/arc-setup-for-login-screen.conf
# Copyright 2017 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Sets up Android directories for the login screen"
author        "chromium-os-dev@chromium.org"

start on start-arc-instance-for-login-screen
stop on stop-arc-instance or stopping ui

env LOGFILE=/var/log/arc-setup-for-login-screen.log
env ENVFILE=/etc/init/arc-setup-env

# Do everything in pre-start to block session_manager. See arc-setup.conf
# for more details.
pre-start script
  bootstat mini-android-start
  exec > $LOGFILE 2>&1
  echo "$(date --rfc-3339=ns): Pre-start arc-setup-for-login-screen"
  set -x
  . $ENVFILE
  /usr/sbin/arc-setup --setup-for-login-screen
  exec bootstat arc-setup-for-mini-android-end
end script

# This file should have neither start nor post-start scripts. See arc-setup.conf
# for more details.

# This script must be in sync with arc-setup.conf's post-stop script.
post-stop script
  exec >> $LOGFILE 2>&1
  echo "$(date --rfc-3339=ns): Post-stop arc-setup-for-login-screen"
  set +e -x
  . $ENVFILE
  stop arc-removable-media
  stop arc-obb-mounter
  stop arc-sdcard
  exec /usr/sbin/arc-setup --stop
end script
