# arc/setup/etc/arc-stop-sysctl.conf
# Copyright 2017 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Tune sysctls for ARC++ stop"
author        "chromium-os-dev@chromium.org"

# See arc-start-sysctl.conf for details.

start on stopped arc-boot-continue
env LOGFILE=/var/log/arc-stop-sysctl.log

script
  {
    echo "Start arc-stop-sysctl"
    KB="$(/usr/share/cros/init/swap.sh get_target_value min_filelist)"
    sysctl vm.min_filelist_kbytes="${KB}"
  } 2>&1 | logger -t "${UPSTART_JOB}"
end script
