# arc/setup/etc/arc-start-sysctl.conf
# Copyright 2017 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Tune sysctls for ARC++ start"
author        "chromium-os-dev@chromium.org"

# vm.min_filelist_kbytes controls how much of the file cache in RAM is
# protected (i.e. not allowed to be evicted).  If too few code pages are kept,
# the system can thrash.  Android apps add pressure to the cache, so we protect
# more pages when ARC++ is running.
#
# See crbug.com/709660 for discussion.
# Also see crbug.com/792152.

start on started arc-boot-continue
env LOGFILE=/var/log/arc-start-sysctl.log

script
  exec > $LOGFILE 2>&1
  echo "$(date --rfc-3339=ns): Start arc-start-sysctl"
  KB="$(/usr/share/cros/init/swap.sh get_target_value min_filelist)"
  sysctl vm.min_filelist_kbytes="${KB}"
end script
