# arc/setup/etc/arc-lifetime.conf
# Copyright 2018 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Clean up the jobs and mount points for the container"
author        "chromium-os-dev@chromium.org"

start on start-arc-instance-for-login-screen
stop on stop-arc-instance or stopping ui

env LOGFILE=/var/log/arc-lifetime.log
env ENVFILE=/etc/init/arc-setup-env

pre-start script
  exec >> $LOGFILE 2>&1
  echo "$(date --rfc-3339=ns): Pre-start arc-lifetime"
end script

# This job does not have script/exec stanza, so it is considered running forever
# once started until stopped.
# See: http://upstart.ubuntu.com/cookbook/#jobs-that-run-forever

post-stop script
  exec > $LOGFILE 2>&1
  echo "$(date --rfc-3339=ns): Post-stop arc-lifetime"
  set +e -x
  . $ENVFILE
  stop arc-removable-media
  stop arc-obb-mounter
  stop arc-sdcard
  exec /usr/sbin/arc-setup --stop
end script
