# arc/setup/etc/arc-boot-continue.conf
# Copyright 2017 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Turn the container into a fully functional one"
author        "chromium-os-dev@chromium.org"

start on continue-arc-boot
stop on stop-arc-instance or stopping ui

env LOGFILE=/var/log/arc-boot-continue.log
env ENVFILE=/etc/init/arc-setup-env

# The following environment variables are passed from session_manager
# and are imported from the event that starts the job.
import ANDROID_DATA_DIR
import ANDROID_DATA_OLD_DIR
import CHROMEOS_DEV_MODE
import CHROMEOS_INSIDE_VM
import CHROMEOS_USER
import CONTAINER_PID
import COPY_PACKAGES_CACHE
import DEMO_SESSION_APPS_PATH
import IS_DEMO_SESSION
import DISABLE_BOOT_COMPLETED_BROADCAST
import ENABLE_VENDOR_PRIVILEGED
import IS_CHILD
import SKIP_PACKAGES_CACHE_SETUP

# Do everything in pre-start to block session_manager.
pre-start script
  bootstat android-start
  exec > $LOGFILE 2>&1
  echo "$(date --rfc-3339=ns): Pre-start arc-boot-continue"
  set -x
  . $ENVFILE
  /usr/sbin/arc-setup --boot-continue
  initctl start -n arc-sdcard CONTAINER_PID=${CONTAINER_PID}
  exec bootstat arc-setup-end
end script

# This file doesn't have a post-stop script. arc-lifetime's cleans things up.
