# init/upstart/test-init/factory-init.conf
description    "Initialize factory test environment"
author         "chromium-os-dev@chromium.org"

# For factory test to perform early initialization Ex: changing UI start
# parameters, disable or add upstart jobs...
# Hint: ui.conf starts on "started boot-services", so this factory-init should
# always run before UI starts.

start on starting boot-services and stopped cros-payloads
task

pre-start script
  . /usr/share/cros/factory_utils.sh
  if ! is_factory_mode; then
    stop
  fi
end script

exec /usr/local/factory/init/startup init
