# init/upstart/test-init/job-filter.conf
# Copyright 2015 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.
description    "Upstart job filter for test environment"
author         "chromium-os-dev@chromium.org"

# This job allows monitoring and filtering Upstart jobs.
# See /usr/sbin/job-filter for more information.

start on starting JOB!=job-filter
task
instance $JOB
exec /usr/sbin/job-filter "$JOB"
