# init/embedded-init/udhcpc.conf
# Copyright (c) 2013 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Run udhcpc dhcp client daemon"
author        "chromium-os-dev@chromium.org"

start on (started boot-services
          and stopped iptables and stopped ip6tables)
stop on stopping boot-services
respawn

pre-start script
  modprobe -q asix || :
  touch /var/run/resolv.conf
end script

expect fork
exec udhcpc
