# imageloader/imageloader.conf
# Copyright 2016 The Chromium OS Authors.All rights reserved.
# Use of this source code is governed by a BSD - style license that can be
# found in the LICENSE file.

description "mount component updated flash if it exists"
author "chromium-os-dev@chromium.org"

start on starting boot-services

# Try to mount flash player if it exists.
# If it does not exist, that is OK: just let imageloader fail.
# This boot time mounting is done because the sandbox/zygote on Chrome
# cannot support loading Flash after their warmup routine. As a performance
# improvement, and architectural simplification, we mount the flash image
# at boot in a pre-determined location that Chrome knows about.
exec /usr/sbin/imageloader --mount \
  --mount_component=PepperFlashPlayer \
  --mount_point=/mnt/stateful_partition/imageloader_mounts/PepperFlashPlayer
