# gobi-cromo-plugin/77-cromo-gobi-device.rules
# do not edit this file, it will be overwritten on update

# If we find a gobi modem let upstart know so it can run cromo
ACTION!="add|change", GOTO="mm_cromo_gobi_device_end"

# If we find a gobi modem let upstart know so it can run cromo
KERNEL=="qcqmi[0-9]*",RUN+="/sbin/initctl emit gobi_modem"

LABEL="mm_cromo_gobi_device_end"

