# gobi-cromo-plugin/77-cromo-gobi-device-blacklist.rules
# do not edit this file, it will be overwritten on update

ACTION!="add|change", GOTO="mm_cromo_gobi_device_blacklist_end"
SUBSYSTEM!="usb", GOTO="mm_cromo_gobi_device_blacklist_end"
ENV{DEVTYPE}!="usb_device",  GOTO="mm_cromo_gobi_device_blacklist_end"
ATTRS{idUsbDriver}!="QCSerial2k", GOTO="mm_cromo_gobi_device_blacklist_end"

# Make sure Modem Manager Ignores all devices handled by QCSerial2k
ENV{ID_MM_DEVICE_IGNORE}="1"

# GOBI Modems are handled by Cromo not by Modem Manager
ATTRS{idVendor}=="0830", ATTRS{idProduct}=="0061", GOTO="mm_cromo_gobi_modem_found"

# Acer Gobi Modem Device
ATTRS{idVendor}=="05c6", ATTRS{idProduct}=="9212", GOTO="mm_cromo_gobi_modem_found"

# HP un2400 Gobi Modem Device
ATTRS{idVendor}=="03f0", ATTRS{idProduct}=="1f1d", GOTO="mm_cromo_gobi_modem_found"

# Panasonic Gobi Modem device
ATTRS{idVendor}=="04da", ATTRS{idProduct}=="250d", GOTO="mm_cromo_gobi_modem_found"

# Dell Gobi Modem device
ATTRS{idVendor}=="413c", ATTRS{idProduct}=="8172", GOTO="mm_cromo_gobi_modem_found"

# Novatel Gobi Modem device
ATTRS{idVendor}=="1410", ATTRS{idProduct}=="a001", GOTO="mm_cromo_gobi_modem_found"

# Asus Gobi Modem device
ATTRS{idVendor}=="0b05", ATTRS{idProduct}=="1776", GOTO="mm_cromo_gobi_modem_found"

# ONDA Gobi Modem device
ATTRS{idVendor}=="19d2", ATTRS{idProduct}=="fff3", GOTO="mm_cromo_gobi_modem_found"

# Generic Gobi Modem device
ATTRS{idVendor}=="05c6", ATTRS{idProduct}=="9001", GOTO="mm_cromo_gobi_modem_found"

# Generic Gobi Modem device
ATTRS{idVendor}=="05c6", ATTRS{idProduct}=="9002", GOTO="mm_cromo_gobi_modem_found"

# Generic Gobi Modem device
ATTRS{idVendor}=="05c6", ATTRS{idProduct}=="9203", GOTO="mm_cromo_gobi_modem_found"

# Generic Gobi Modem device
ATTRS{idVendor}=="05c6", ATTRS{idProduct}=="9222", GOTO="mm_cromo_gobi_modem_found"

# Dell Gobi 2000 Modem device (N0218, VU936)
ATTRS{idVendor}=="413c", ATTRS{idProduct}=="8186", GOTO="mm_cromo_gobi_modem_found"

# Sony Gobi 2000 Modem device (N0279, VU730)
ATTRS{idVendor}=="05c6", ATTRS{idProduct}=="9225", GOTO="mm_cromo_gobi_modem_found"

# Samsung Gobi 2000 Modem device (VL176)
ATTRS{idVendor}=="05c6", ATTRS{idProduct}=="9245", GOTO="mm_cromo_gobi_modem_found"

# HP Gobi 2000 Modem device (VP412)
ATTRS{idVendor}=="03f0", ATTRS{idProduct}=="251d", GOTO="mm_cromo_gobi_modem_found"

# Acer Gobi 2000 Modem device (VP413)
ATTRS{idVendor}=="05c6", ATTRS{idProduct}=="9215", GOTO="mm_cromo_gobi_modem_found"

# Asus Gobi 2000 Modem device (VR305)
ATTRS{idVendor}=="05c6", ATTRS{idProduct}=="9265", GOTO="mm_cromo_gobi_modem_found"

# Top Global Gobi 2000 Modem device (VR306)
ATTRS{idVendor}=="05c6", ATTRS{idProduct}=="9235", GOTO="mm_cromo_gobi_modem_found"

# iRex Technologies Gobi 2000 Modem device (VR307)
ATTRS{idVendor}=="05c6", ATTRS{idProduct}=="9275", GOTO="mm_cromo_gobi_modem_found"

# Sierra Wireless Gobi 2000 Modem device (VT773)
ATTRS{idVendor}=="1199", ATTRS{idProduct}=="9001", GOTO="mm_cromo_gobi_modem_found"

# Sierra Wireless Gobi 2000 Modem device (VT773)
ATTRS{idVendor}=="1199", ATTRS{idProduct}=="9002", GOTO="mm_cromo_gobi_modem_found"

# Sierra Wireless Gobi 2000 Modem device (VT773)
ATTRS{idVendor}=="1199", ATTRS{idProduct}=="9003", GOTO="mm_cromo_gobi_modem_found"

# Sierra Wireless Gobi 2000 Modem device (VT773)
ATTRS{idVendor}=="1199", ATTRS{idProduct}=="9004", GOTO="mm_cromo_gobi_modem_found"

# Sierra Wireless Gobi 2000 Modem device (VT773)
ATTRS{idVendor}=="1199", ATTRS{idProduct}=="9005", GOTO="mm_cromo_gobi_modem_found"

# Sierra Wireless Gobi 2000 Modem device (VT773)
ATTRS{idVendor}=="1199", ATTRS{idProduct}=="9006", GOTO="mm_cromo_gobi_modem_found"

# Sierra Wireless Gobi 2000 Modem device (VT773)
ATTRS{idVendor}=="1199", ATTRS{idProduct}=="9007", GOTO="mm_cromo_gobi_modem_found"

# Sierra Wireless Gobi 2000 Modem device (VT773)
ATTRS{idVendor}=="1199", ATTRS{idProduct}=="9008", GOTO="mm_cromo_gobi_modem_found"

# Sierra Wireless Gobi 2000 Modem device (VT773)
ATTRS{idVendor}=="1199", ATTRS{idProduct}=="9009", GOTO="mm_cromo_gobi_modem_found"

# Sierra Wireless Gobi 2000 Modem device (VT773)
ATTRS{idVendor}=="1199", ATTRS{idProduct}=="900a", GOTO="mm_cromo_gobi_modem_found"

# CMDTech Gobi 2000 Modem device (VU922)
ATTRS{idVendor}=="16d8", ATTRS{idProduct}=="8002", GOTO="mm_cromo_gobi_modem_found"

# Gobi 2000 Modem device
ATTRS{idVendor}=="05c6", ATTRS{idProduct}=="9205", GOTO="mm_cromo_gobi_modem_found"

# Nothing matched
GOTO="mm_cromo_gobi_device_blacklist_end"

#
# If we found a modem device let upstart know so it can run cromo
#
LABEL="mm_cromo_gobi_modem_found"
RUN+="/sbin/initctl emit gobi_modem"

LABEL="mm_cromo_gobi_device_blacklist_end"

