# gobi-cromo-plugin/record-modem.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# Record modem data to fill in our database of MEIDs for about 200
# dogfood machines where we have only IMEIs

start on login-prompt-ready

task

script
        /usr/bin/record_modem
end script
