# gobi-cromo-plugin/usbsrvd.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# usbservd
#
# starts usbsrvd if a GOBI modem has been detected

start on gobi_modem

expect fork

# Set up all necessary configuratino
pre-start script
  if [ -x /usr/local/bin/usbsrvd ] ; then
     modprobe tusbd
     iptables -A INPUT -p tcp --dport 32032 -j ACCEPT
  fi
end script


post-start script
  deviceid=$(/usr/local/bin/usbsrv -l | sed -n -e /Qualcomm/s/:.*//p | head -1)
  /usr/local/bin/usbsrv -s $deviceid
end script

# enable usbsrvd
exec /usr/local/bin/usbsrvd

