# crash-reporter/init/crash-sender.service
[Unit]
Description=Crash Sender

[Service]
ExecStart=/sbin/crash_sender
