# crash-reporter/init/crash-reporter.service
[Unit]
Description=Initialize Crash Reporter

[Service]
Type=oneshot
ExecStart=/sbin/crash_reporter --init

[Install]
WantedBy=multi-user.target
