# app_shell_launcher/init/boot-complete.conf
# Copyright 2014 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Boot marker signaling that system boot is complete"
author        "chromium-os-dev@chromium.org"

# TODO(derat): Defer starting until app_shell has loaded the app, maybe.
# This is needed to run permission_broker, and we don't know how soon the app
# will need it, though...
start on started app_shell
