# peerd/init/peerd.conf
# Copyright 2014 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Run the peer discovery services daemon"
author        "chromium-os-dev@chromium.org"

start on starting system-services
stop on stopping system-services

env PEERD_LOG_LEVEL=0

# Minijail actually forks off our desired process
expect fork

exec minijail0 -i -g peerd -u peerd \
    /usr/bin/peerd --log-level="${PEERD_LOG_LEVEL}"
