# hammerd/init/hammerd.conf
# Copyright 2017 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Hammer firmware updater task"
author        "chromium-os-dev@chromium.org"

start on hammer-device-added

# minijail is going to fork hammerd.
expect fork
task

# Trigger hammerd process.
# Note: minijail0 "-N" argument (cgroup namespace) is only supported at
# kernel >= 3.14. Remove it if the hammerd runs at the kernel < 3.14.
exec /sbin/minijail0 -e -N -p -l -u hammerd -g hammerd -i -c 0002 \
     /usr/bin/hammerd
