# power_manager/udev/99-input-wakeup-powerd-chown.rules
ACTION!="remove", SUBSYSTEM=="input", ATTRS{power/wakeup}=="*", RUN+="/bin/chown power:power /sys/class/input/%k/device/device/power/wakeup"
