# power_manager/udev/99-usb-hid-wake.rules
ACTION!="remove", SUBSYSTEM=="usb", RUN+="/lib/udev/usb-hid-wake.sh %k"

