# chaps/chapsd.conf
# Copyright (c) 2011 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Start the Chaps daemon (chapsd)"
author          "chromium-os-dev@chromium.org"

# Starts the Chaps daemon, which provides TPM-backed cryptographic services.

start on (started tcsd and started dbus)
stop on (starting reboot or starting halt)
respawn

script
  VERBOSE=""
  if [ -f /home/chronos/.chaps_debug_2 ]; then
    VERBOSE="--v=2"
    rm /home/chronos/.chaps_debug_2
  elif [ -f /home/chronos/.chaps_debug_1 ]; then
    VERBOSE="--v=1"
    rm /home/chronos/.chaps_debug_1
  fi
  exec chapsd ${VERBOSE}
end script

