# chaps/chapsd.conf
# Copyright (c) 2011 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Start the Chaps daemon (chapsd)"
author          "chromium-os-dev@chromium.org"

# Starts the Chaps daemon, which provides TPM-backed cryptographic services.

start on (started tcsd and started dbus)
stop on (starting reboot or starting halt)
respawn

script
  LIB=""
  if [ -f /home/chronos/.disable_chaps ]; then
    # Use the openCryptoki PKCS #11 implementation.
    LIB=--lib=libopencryptoki.so
  fi
  exec chapsd ${LIB}
end script

