# init/upstart/file_attrs_cleaner_tool.conf
# Copyright 2019 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Run file_attrs_cleaner_tool on the encrypted home directories."
author          "chromium-os-dev@chromium.org"

start on started system-services

# See crbug.com/915846 and crbug.com/919486.
exec ionice -c 3 -- /sbin/file_attrs_cleaner_tool /mnt/stateful_partition/home
