# init/upstart/pstore.conf
# Copyright (c) 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "mount pstore on /dev/pstore"
author          "chromium-os-dev@chromium.org"

start on starting boot-services
task

script
  mkdir -p /dev/pstore
  mount -n -t pstore -o nodev,nosuid,noexec none /dev/pstore
  chmod 0710 /dev/pstore
  # TODO(ellyjones): restrict more tightly (crosbug.com/35122)
  chgrp debugd /dev/pstore
end script
