# init/upstart/pstore.conf
# Copyright (c) 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "mount pstore for kernel buffer over reboots"
author          "chromium-os-dev@chromium.org"

start on starting boot-services
task

script
  path="/sys/fs/pstore"
  mount -n -t pstore -o nodev,nosuid,noexec pstore "${path}"
  chmod 0750 "${path}"
  # TODO(ellyjones): restrict more tightly (crosbug.com/35122)
  chgrp debugd "${path}"

  # Legacy path support.  Can probably drop once everything uses sysfs.
  mkdir -p /dev/pstore
  mount --bind "${path}" /dev/pstore
end script
