# init/upstart/static-nodes.conf
# Copyright 2016 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Create the registered static device nodes"
author          "chromium-os-dev@chromium.org"

start on starting boot-services
task

# Create static device nodes from modules.devname.
exec static_node_tool
