# cryptohome/init/cryptohomed-client.conf
# Copyright (c) 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Stop existing cryptohomed clients"
author          "chromium-os-dev@chromium.org"

# Any long-running invocations of the client must first be terminated.
start on stopping cryptohomed
task
exec pkill -TERM cryptohome
