# apmanager/init/apmanager.conf
# Copyright 2014 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Run the access point manager daemon"
author        "chromium-os-dev@chromium.org"

start on starting system-services
stop on stopping system-services

env APMANAGER_LOG_LEVEL=0

exec /usr/bin/apmanager --v="${APMANAGER_LOG_LEVEL}"
