# apmanager/init/apmanager.conf
# Copyright 2014 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Run the access point manager daemon"
author        "chromium-os-dev@chromium.org"

start on starting system-services
stop on stopping system-services
expect fork

env APMANAGER_LOG_LEVEL=0

pre-start script
  # Load the module that provides the WiFi configuration API, since
  # apmanager will abort if that API is not available. In most cases,
  # cfg80211 will be loaded implicitly when the device driver is
  # loaded (in preload-network).  However, this deals with the
  # first-boot case, in case apmanager starts before the device driver is
  # loaded.
  modprobe cfg80211 ||
    logger -p err -t "$UPSTART_JOB" "Failed to load cfg80211"

  # Create directory for storing config files.
  mkdir -m 0755 -p /var/run/apmanager/hostapd
  mkdir -m 0755 -p /var/run/apmanager/dnsmasq
  chown -R apmanager:apmanager /var/run/apmanager/hostapd
  chown -R apmanager:apmanager /var/run/apmanager/dnsmasq
end script

exec /usr/bin/apmanager --v="${APMANAGER_LOG_LEVEL}"
