# camera/common/init/cros-camera-gpu-algo.conf
# Copyright 2019 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Start camera GPU algorithm service"
author        "chromium-os-dev@chromium.org"

start on started system-services
stop on stopping system-services
expect fork
respawn
respawn limit 10 5

script
  # TODO(hywu): minimum sandboxing with minijail
  exec sudo -u arc-camera -g arc-camera /usr/bin/cros_camera_algo --type=gpu
end script
