# camera/common/init/cros-camera-gpu-algo.conf
# Copyright 2019 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Start camera GPU algorithm service"
author        "chromium-os-dev@chromium.org"

start on starting cros-camera
stop on stopping cros-camera
oom score 600
expect fork
respawn
respawn limit 10 5

env SECCOMP_POLICY_FILE=/usr/share/policy/cros-camera-gpu-algo.policy

script
  # Start constructing minijail0 args...
  set --

  # Enter a new mount, network, PID, IPC and cgroup namespace.
  set -- "$@" -v -e -p -l -N

  # Change user and group to arc-camera.
  set -- "$@" -u arc-camera -g arc-camera

  # Set -i to fork and daemonize an init-like process that Upstart will track
  # as the service.
  set -- "$@" -i

  # Chroot and mount /dev, /proc, /tmp and /run/camera. This service uses
  # /dev/urandom and /dev/mali0.
  set -- "$@" -P /mnt/empty -b / -b /proc -b /dev,,1
  set -- "$@" -k 'tmpfs,/run,tmpfs,MS_NOSUID|MS_NODEV|MS_NOEXEC'
  set -- "$@" -b /run/camera,,1 -b /run/systemd/journal -t

  # Assume static ELF binary to give arc-camera access to /proc/self.
  set -- "$@" -Tstatic

  # Drop privileges and set seccomp filter.
  set -- "$@" -n -S "${SECCOMP_POLICY_FILE}"
  set -- "$@" -- /usr/bin/cros_camera_algo --type=gpu

  exec minijail0 "$@"
end script
