# timberslide/init/timberslide.conf
# Copyright 2015 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "EC log concatenator"
author        "chromium-os-dev@chromium.org"

start on started system-services
stop on stopping system-services
# 69 is EX_UNAVAILABLE which we return if there is no EC log to pull.
normal exit 0 69
respawn
respawn limit 3 10  # if the job respawns 3 times in 10 seconds, stop trying.

exec timberslide >/var/log/timberslide.out 2>&1
