# chaps/init/chapsd.conf
# Copyright (c) 2011 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Start the Chaps daemon (chapsd)"
author          "chromium-os-dev@chromium.org"

# Starts the Chaps daemon, which provides TPM-backed cryptographic services.

start on started tcsd and started boot-services
stop on stopping boot-services
respawn

pre-start script
  mkdir -p /var/lib/chaps
  chown chaps:chronos-access /var/lib/chaps
end script

exec /usr/sbin/chapsd --auto_load_system_token
