# arc/removable-media/arc-removable-media-write.conf
# Copyright 2018 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Mount /media/removable for ARC."
author        "chromium-os-dev@chromium.org"

start on started arc-network
stop on stopping arc-network

env SOURCE=/media/removable
env DEST=/run/arc/media/removable-write
env UMASK=007
env ANDROID_ROOT_UID=0
env ANDROID_EVERYBODY_GID=9997

script
  {
    echo "Start arc-removable-media-write"
    set -x
    # Unlike other services, redirect to logger, since libfuse used in
    # mount-passthrough uses stdio for logging.
    exec mount-passthrough-jailed "${SOURCE}" "${DEST}" "${UMASK}" \
        "${ANDROID_ROOT_UID}" "${ANDROID_EVERYBODY_GID}" "write"
  } 2>&1 | logger -t "${UPSTART_JOB}"
end script

post-stop exec logger -t "${UPSTART_JOB}" "Post-stop arc-removable-media-write"
