# init/upstart/workaround-init/midi-workaround.conf
# Copyright 2016 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Workaround for kernel MIDI driver bugs"
author          "chromium-os-dev@chromium.org"

# This needs to finish before Chrome launches.
start on starting ui
task

# Workaround for MIDI devices not being detected if none have been
# connected before loading a Web MIDI page.
# See https://crbug.com/499817.
# This was fixed in linux-3.19+.

script
  modprobe -q snd_seq_midi \
    || logger -t "$UPSTART_JOB" "modprobe snd_seq_midi failed."
end script
