# ippusb_manager/etc/init/ippusb-pre-upstart-socket-bridge.conf
# Copyright 2017 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description "Prepare upstart-socket-bridge for ippusb"

start on starting upstart-socket-bridge

task

# HACK: Create the directory now, since upstart-socket-bridge doesn't create it.
script
  ippusb_socket_dir="/run/ippusb/"
  ippusb_group="ippusb"

  mkdir -p -m 0775 "${ippusb_socket_dir}"
  chmod 0775 "${ippusb_socket_dir}"
  chown root:"${ippusb_group}" "${ippusb_socket_dir}"
end script
