# attestation/server/attestationd.conf
# Copyright 2014 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Chromium OS device attestation service."
author          "chromium-os-dev@chromium.org"

start on started tcsd and started boot-services
stop on stopping boot-services
respawn

exec /usr/sbin/attestationd
