# attestation/server/attestationd.conf
# Copyright 2014 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Chromium OS device attestation service."
author          "chromium-os-dev@chromium.org"

start on starting system-services
stop on stopping system-services
expect fork
respawn

# TODO(namnguyen): further restrict privileges, see crbug.com/357290.
exec minijail0 -i -n -u attestation -g attestation /usr/sbin/attestationd
